# testscript.service
[Unit]
Description=Test service for persistence

[Service]
ExecStart=/usr/bin/testscript.sh

[Install]
WantedBy=multi-user.target

